# guest/openspan.service
[Unit]
Description=OpenSpan HID bridge daemon
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/bin/sh -c 'for i in 1 2 3 4 5 6 7 8 9 10; do hciconfig hci0 up 2>/dev/null && hciconfig hci0 class 0x0005C0 2>/dev/null && break; sleep 1; done'
ExecStart=/usr/bin/python3 /opt/openspan/openspan_hid.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
