# guest/openspan-wireplumber.service
[Unit]
Description=OpenSpan WirePlumber (system)
After=openspan-pipewire.service openspan-dbus.service bluetooth.service
Requires=openspan-pipewire.service openspan-dbus.service
Wants=bluetooth.service

[Service]
Type=simple
Environment=XDG_RUNTIME_DIR=/run/openspan
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/openspan/bus
# Wait for BOTH the pipewire socket and the session bus to exist -- no race.
ExecStartPre=/bin/sh -c 'for i in $(seq 40); do [ -S /run/openspan/pipewire-0 ] && [ -S /run/openspan/bus ] && exit 0; sleep 0.5; done'
ExecStart=/usr/bin/wireplumber
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
