# guest/system/openspanble.service.d/10-wait.conf
[Service]
# Reset the inherited ExecStartPre list, then wait for the radio and set
# dual-mode ONLY IF NEEDED (conditional, via ensure-dualmode.sh) so restarting
# the keyboard daemon never power-cycles the radio and drops active A2DP audio.
ExecStartPre=
ExecStartPre=/opt/openspan/wait-hci0.sh
ExecStartPre=/opt/openspan/ensure-dualmode.sh
