# guest/system/openspanble.service.d/override.conf
[Unit]
After=openspan-btready.service openspan-wireplumber.service bluetooth.service
Wants=bluetooth.service
[Service]
Restart=always
RestartSec=2
